# installer/barprep-edge.service
[Unit]
Description=BarPrep Edge OS Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/barprep-edge.env
ExecStart=/opt/barprep-edge/.venv/bin/barprep-edge
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
